# server/caddy.conf
[program:caddy]
command=caddy -conf=/root/mediasoup-demo/server/Caddyfile -http2=false
process_name=%(program_name)s
numprocs=1
directory=/root/mediasoup-demo/server
autostart=true
user=root
